# system.d/service.webgrabplus.timer
[Unit]
Description=WebGrab+Plus daily

[Timer]
OnCalendar=*-*-* 03:02:01
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target